# scripts/pushWeather.timer
[Unit]
Description=push weather data to django

[Timer]
OnBootSec=0min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
